# systemd/gpu-cooling.service
[Unit]
Description=Boot-time GPU cooling profile (persistence, 260 W cap, chassis fans, GPU fan floor)
Documentation=file:///%h/AI/2x-3090-GA102-300-A1-sglang-inference/README.md
# Order after the NVIDIA driver/persistence stack (the actual dependency), NOT
# after multi-user.target — the latter closed a dependency cycle with
# gpu-fan-curve + swebench-bakeoff (both WantedBy multi-user.target), so systemd
# silently dropped swebench-bakeoff at every boot and the bake-off never
# auto-resumed after a reboot. nvidia-persistenced is outside that loop.
After=nvidia-persistenced.service
Wants=nvidia-persistenced.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/gpu-cooling.sh
# Writes /sys/class/hwmon and runs nvidia-smi -pm/-pl; needs root.
User=root
Group=root

[Install]
WantedBy=multi-user.target
